# lib64/apk/packages/dawn.list
/etc/config/dawn
/etc/init.d/dawn
/usr/sbin/dawn
