# lib64/apk/packages/https-dns-proxy.list
/etc/config/https-dns-proxy
/etc/init.d/https-dns-proxy
/etc/uci-defaults/50-https-dns-proxy-migrate-options.sh
/usr/sbin/https-dns-proxy
