# lib64/apk/packages/zapret.list
/etc/hotplug.d/iface/90-zapret
/etc/init.d/zapret
/etc/uci-defaults/zapret-uci-def-cfg.sh
/opt/zapret/blockcheck.sh
/opt/zapret/comfunc.sh
/opt/zapret/common/base.sh
/opt/zapret/common/custom.sh
/opt/zapret/common/dialog.sh
/opt/zapret/common/elevate.sh
/opt/zapret/common/fwtype.sh
/opt/zapret/common/installer.sh
/opt/zapret/common/ipt.sh
/opt/zapret/common/linux_daemons.sh
/opt/zapret/common/linux_fw.sh
/opt/zapret/common/linux_iphelper.sh
/opt/zapret/common/list.sh
/opt/zapret/common/nft.sh
/opt/zapret/common/pf.sh
/opt/zapret/common/virt.sh
/opt/zapret/config
/opt/zapret/config.default
/opt/zapret/def-cfg.sh
/opt/zapret/dwc.sh
/opt/zapret/files/fake/4pda.bin
/opt/zapret/files/fake/dht_find_node.bin
/opt/zapret/files/fake/dht_get_peers.bin
/opt/zapret/files/fake/discord-ip-discovery-with-port.bin
/opt/zapret/files/fake/discord-ip-discovery-without-port.bin
/opt/zapret/files/fake/dtls_clienthello_w3_org.bin
/opt/zapret/files/fake/http_iana_org.bin
/opt/zapret/files/fake/isakmp_initiator_request.bin
/opt/zapret/files/fake/max.bin
/opt/zapret/files/fake/quic_initial_facebook_com.bin
/opt/zapret/files/fake/quic_initial_facebook_com_quiche.bin
/opt/zapret/files/fake/quic_initial_rr1---sn-xguxaxjvh-n8me_googlevideo_com_kyber_1.bin
/opt/zapret/files/fake/quic_initial_rr1---sn-xguxaxjvh-n8me_googlevideo_com_kyber_2.bin
/opt/zapret/files/fake/quic_initial_rr2---sn-gvnuxaxjvh-o8ge_googlevideo_com.bin
/opt/zapret/files/fake/quic_initial_rutracker_org.bin
/opt/zapret/files/fake/quic_initial_rutracker_org_kyber_1.bin
/opt/zapret/files/fake/quic_initial_rutracker_org_kyber_2.bin
/opt/zapret/files/fake/quic_initial_vk_com.bin
/opt/zapret/files/fake/quic_initial_www_google_com.bin
/opt/zapret/files/fake/quic_short_header.bin
/opt/zapret/files/fake/stun.bin
/opt/zapret/files/fake/t2.bin
/opt/zapret/files/fake/tls_clienthello_gosuslugi_ru.bin
/opt/zapret/files/fake/tls_clienthello_iana_org.bin
/opt/zapret/files/fake/tls_clienthello_max_ru.bin
/opt/zapret/files/fake/tls_clienthello_rutracker_org_kyber.bin
/opt/zapret/files/fake/tls_clienthello_sberbank_ru.bin
/opt/zapret/files/fake/tls_clienthello_vk_com.bin
/opt/zapret/files/fake/tls_clienthello_vk_com_kyber.bin
/opt/zapret/files/fake/tls_clienthello_www_google_com.bin
/opt/zapret/files/fake/tls_clienthello_www_onetrust_com.bin
/opt/zapret/files/fake/wireguard_initiation.bin
/opt/zapret/files/fake/wireguard_response.bin
/opt/zapret/files/fake/zero_1024.bin
/opt/zapret/files/fake/zero_256.bin
/opt/zapret/files/fake/zero_512.bin
/opt/zapret/files/huawei/E8372/run-zapret-hostlist
/opt/zapret/files/huawei/E8372/run-zapret-ip
/opt/zapret/files/huawei/E8372/unfuck_nfqueue.ko
/opt/zapret/files/huawei/E8372/unzapret
/opt/zapret/files/huawei/E8372/unzapret-ip
/opt/zapret/files/huawei/E8372/zapret
/opt/zapret/files/huawei/E8372/zapret-ip
/opt/zapret/init.d/openwrt/90-zapret
/opt/zapret/init.d/openwrt/custom.d/.keep
/opt/zapret/init.d/openwrt/custom.d/10-script.sh
/opt/zapret/init.d/openwrt/custom.d/20-script.sh
/opt/zapret/init.d/openwrt/custom.d/50-script.sh
/opt/zapret/init.d/openwrt/custom.d/60-script.sh
/opt/zapret/init.d/openwrt/custom.d/90-script.sh
/opt/zapret/init.d/openwrt/firewall.zapret
/opt/zapret/init.d/openwrt/functions
/opt/zapret/init.d/openwrt/zapret
/opt/zapret/ip2net/ip2net
/opt/zapret/ipset/antifilter.helper
/opt/zapret/ipset/clear_lists.sh
/opt/zapret/ipset/create_ipset.sh
/opt/zapret/ipset/cust1.txt
/opt/zapret/ipset/cust2.txt
/opt/zapret/ipset/cust3.txt
/opt/zapret/ipset/cust4.txt
/opt/zapret/ipset/def.sh
/opt/zapret/ipset/get_antifilter_allyouneed.sh
/opt/zapret/ipset/get_antifilter_ip.sh
/opt/zapret/ipset/get_antifilter_ipresolve.sh
/opt/zapret/ipset/get_antifilter_ipsmart.sh
/opt/zapret/ipset/get_antifilter_ipsum.sh
/opt/zapret/ipset/get_antizapret_domains.sh
/opt/zapret/ipset/get_config.sh
/opt/zapret/ipset/get_exclude.sh
/opt/zapret/ipset/get_ipban.sh
/opt/zapret/ipset/get_reestr_preresolved.sh
/opt/zapret/ipset/get_reestr_preresolved_smart.sh
/opt/zapret/ipset/get_reestr_resolvable_domains.sh
/opt/zapret/ipset/get_refilter_domains.sh
/opt/zapret/ipset/get_refilter_ipsum.sh
/opt/zapret/ipset/get_user.sh
/opt/zapret/ipset/zapret-hosts-auto.txt
/opt/zapret/ipset/zapret-hosts-google.txt
/opt/zapret/ipset/zapret-hosts-user-exclude.txt
/opt/zapret/ipset/zapret-hosts-user-exclude.txt.default
/opt/zapret/ipset/zapret-hosts-user.txt
/opt/zapret/ipset/zapret-ip-exclude.txt
/opt/zapret/ipset_def/zapret-hosts-google.txt
/opt/zapret/ipset_def/zapret-hosts-user-exclude.txt
/opt/zapret/ipset_def/zapret-hosts-user.txt
/opt/zapret/ipset_def/zapret-ip-exclude.txt
/opt/zapret/mdig/mdig
/opt/zapret/nfq/nfqws
/opt/zapret/renew-cfg.sh
/opt/zapret/restore-def-cfg.sh
/opt/zapret/script-exec.sh
/opt/zapret/sync_config.sh
/opt/zapret/uci-def-cfg.sh
/opt/zapret/update-pkg.sh
